# django-jenkins-for-testing/roles/mysql/templates/create-mysql-tmpfs.conf
start on local-filesystems

console log

script
    cp -Rfa /var/lib/mysql /mnt/mysql_tmpfs
    service mysql start
end script
